# packages/d/dhcpcd/files/dhcpcd.service
[Unit]
Description=DHCP Client
After=network.target
[Service]
Type=forking
ExecStart=/sbin/dhcpcd
ExecStop=/sbin/dhcpcd -k
[Install]
WantedBy=multi-user.target

