# packages/d/dhcpcd/files/dhcpcd.service
[Unit]
Description=DHCP Client
After=network.target
[Service]
Type=forking
ExecStart=/usr/sbin/dhcpcd
ExecStop=/usr/sbin/dhcpcd -k
[Install]
WantedBy=multi-user.target

