# packages/o/openssh/files/sshd@.service
[Unit]
Description=OpenSSH per-connection server daemon
Documentation=man:sshd(8) man:sshd_config(5)
Wants=sshd-keygen.target
After=sshd-keygen.target syslog.target auditd.service

[Service]
ExecStart=-/usr/sbin/sshd -i -e $OPTIONS
StandardInput=socket
StandardError=syslog
