# packages/v/virtualbox/files/vboxservice.service
[Unit]
Description=VirtualBox guest service
ConditionVirtualization=oracle

[Service]
ExecStart=/usr/sbin/VBoxService -f

[Install]
WantedBy=multi-user.target
