# packages/v/virtualbox/files/vboxdrv.service
[Unit]
SourcePath=/usr/lib64/virtualbox/vboxdrv.sh
Description=VirtualBox Linux kernel module
Before=runlevel2.target runlevel3.target runlevel4.target runlevel5.target shutdown.target
After=
Conflicts=shutdown.target

[Service]
Type=forking
Restart=no
TimeoutSec=5min
IgnoreSIGPIPE=no
KillMode=process
GuessMainPID=no
RemainAfterExit=yes
ExecStart=/usr/lib64/virtualbox/vboxdrv.sh start
ExecStop=/usr/lib64/virtualbox/vboxdrv.sh stop

[Install]
WantedBy=multi-user.target
