# packages/f/flatpak/files/user/flatpak-update.service
[Unit]
Description=Update Flatpaks

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --assumeyes --noninteractive --user
ExecStart=/usr/bin/flatpak remove --assumeyes --noninteractive --user --unused
