# packages/d/dnsmasq/files/dnsmasq.service
[Unit]
Description=DNS cache and forwarder
After=network.target

[Service]
Type=dbus
BusName=uk.org.thekelleys.dnsmasq
ExecStartPre=/usr/bin/dnsmasq --test
ExecStart=/usr/bin/dnsmasq -k --enable-dbus
ExecReload=/bin/kill -HUP $MAINPID

[Install]
WantedBy=multi-user.target
