# packages/f/flatpak/files/flatpak-enable-flathub.service
[Unit]
Description=Add flathub flatpak repositories
ConditionPathExists=!/var/lib/flatpak/.flathub-initialized
Before=flatpak-system-helper.service
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/flatpak remote-add --system --if-not-exists flathub https://dl.flathub.org/repo/flathub.flatpakrepo
ExecStartPost=/usr/bin/touch /var/lib/flatpak/.flathub-initialized

[Install]
WantedBy=multi-user.target
