# packages/l/libvirt/files/virt-secret-init-encryption.service
[Service]
Type=oneshot
ExecStartPre=/usr/bin/sh -c 'umask 0077 && /usr/bin/mkdir -p /var/lib/libvirt/secrets/'
