# packages/u/usysconf-epoch/files/epoch.service
[Unit]
Description=Perform the epoch transition
After=network-online.target
Wants=network-online.target

[Service]
EnvironmentFile=-/etc/sysconfig/epoch
ExecStart=/usr/lib64/usysconf/epoch.sh
Type=oneshot

RemainAfterExit=yes
TimeoutStartSec=60s

[Install]
WantedBy=multi-user.target
