# packages/u/usysconf-epoch/files/usysconf-epoch.service
[Unit]
Description=Usysconf Epoch script

DefaultDependencies=no
After=local-fs.target systemd-tmpfiles-setup.service
Before=sysinit.target
Conflicts=shutdown.target initrd-switch-root.target
Before=shutdown.target initrd-switch-root.target

[Service]
ExecStart=/usr/lib64/usysconf/epoch.sh
Type=oneshot
RemainAfterExit=yes

[Install]
WantedBy=sysinit.target
