# packages/o/openrgb/files/openrgb.service
[Unit]
Description=OpenRGB daemon

[Service]
Type=simple
ExecStart=/usr/bin/openrgb --server
Restart=on-failure

[Install]
WantedBy=multi-user.target
