# packages/b/bzflag/files/bzflag.service
[Unit]
Description=BZFlag game server
Documentation=man:bzfs(6)
After=network.target

[Service]
Restart=on-failure
User=bzflag
Group=games
ExecStart=/usr/bin/bzfs -advertise NONE

[Install]
WantedBy=multi-user.target
