# packages/u/unbound/files/unbound.service
[Unit]
Description=Unbound recursive Domain Name Server
Documentation=man:unbound(8)
After=network.target
Before=nss-lookup.target
Wants=nss-lookup.target

[Service]
ExecReload=/bin/kill -HUP $MAINPID
ExecStart=/usr/sbin/unbound -d
Restart=always

[Install]
WantedBy=multi-user.target
