# packages/m/minidlna/files/minidlna.service
[Unit]
Description=minidlna server
After=network.target

[Service]
Type=simple
User=minidlna
Group=minidlna
ExecStart=/usrs/bin/minidlnad -S

[Install]
WantedBy=multi-user.target
