# packages/a/apcupsd/files/apcupsd.service
[Unit]
Description=APC UPS Power Control Daemon for Linux
After=network-online.target
Wants=network-online.target

[Service]
ExecStartPre=-/bin/rm -f /etc/apcupsd/powerfail
ExecStart=/usr/sbin/apcupsd -b
KillMode=process

[Install]
WantedBy=multi-user.target
