# packages/f/flatpak/files/system/flatpak-update.service
[Unit]
Description=Update Flatpaks

[Service]
Type=oneshot
ExecStart=/usr/bin/flatpak update --assumeyes --noninteractive --system
ExecStart=/usr/bin/flatpak remove --assumeyes --noninteractive --system --unused
