# packages/f/flatpak/files/system/flatpak-add-flathub-repos.service
[Unit]
Description=Adds Flathub flatpak repositories
ConditionPathExists=!/var/lib/flatpak/.solus-initialized
Before=flatpak-system-helper.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/bin/flatpak remote-add --system --if-not-exists flathub https://dl.flathub.org/repo/flathub.flatpakrepo
ExecStartPost=/usr/bin/touch /var/lib/flatpak/.solus-initialized

[Install]
WantedBy=multi-user.target