# packages/w/waydroid/files/waydroid-pre.service
[Unit]
Description=Steps that should run before Waydroid starts
Wants=dev-binderfs.mount

[Service]
ExecStart=/usr/bin/ln -sf /dev/binderfs/binder /dev/binderfs/vndbinder /dev/binderfs/hwbinder /dev/
Type=oneshot
RemainAfterExit=yes
