# packages/m/mongodb/files/mongodb.service
[Unit]
Description=High-performance, schema-free document-oriented database
After=network.target

[Service]
User=mongodb
Group=mongodb
PermissionsStartOnly=true
ExecStart=/usr/share/mongodb/mongodb.sh
PIDFile=/run/mongodb/mongod.pid
PrivateTmp=true

[Install]
WantedBy=multi-user.target
