# packages/n/netdata/files/netdata-setup.service
[Unit]
Description=Netdata setup service
ConditionPathExists=|!/etc/netdata/netdata.conf

[Service]
ExecStart=/usr/bin/touch /etc/netdata/netdata.conf
ExecStartPost=/usr/bin/chown netdata:netdata /etc/netdata/netdata.conf
Type=oneshot
RemainAfterExit=yes
