# packages/l/lightdm/files/lightdm.service
[Unit]
Description=Display Manager
Conflicts=xdm.service gdm.service kdm.service lxdm.service slim.service
After=systemd-user-sessions.service getty@tty7.service plymouth-quit.service systemd-logind.service

[Service]
ExecStart=/usr/sbin/lightdm
Restart=always
IgnoreSIGPIPE=no
BusName=org.freedesktop.DisplayManager

[Install]
Alias=displaymanager.service
WantedBy=graphical.target

