# packages/l/lightdm/files/lightdm.service
[Unit]
Description=Display Manager
Conflicts=xdm.service gdm.service kdm.service lxdm.service slim.service

[Service]
ExecStart=/usr/sbin/lightdm

[Install]
Alias=displaymanager.service
WantedBy=graphical.target

