# packages/d/dnscrypt-proxy/files/dnscrypt-proxy.service
[Unit]
Description=Encrypted/authenticated DNS proxy
ConditionFileIsExecutable=/usr/bin/dnscrypt-proxy

[Service]
StartLimitInterval=5
StartLimitBurst=10
ExecStart=/usr/bin/dnscrypt-proxy --config /etc/dnscrypt-proxy.toml 

Restart=always
RestartSec=120
EnvironmentFile=-/etc/sysconfig/dnscrypt-proxy

[Install]
WantedBy=multi-user.target
