# packages/x/xdg-user-dirs/files/xdg-user-dirs.service
[Unit]
Description=Initialise XDG User Directories

[Service]
Type=oneshot
ExecStart=/usr/bin/xdg-user-dirs-update

[Install]
WantedBy=default.target
