# packages/i/ivpn/files/ivpn-service.service
[Unit]
Description=ivpn-service

[Service]
Type=simple
User=root
Group=root
ExecStart=/usr/bin/ivpn-service
Restart=always
WorkingDirectory=/
TimeoutStopSec=infinity

[Install]
WantedBy=multi-user.target